# .config/systemd/user/oidc-agent.service
[Unit]
Description=OIDC Agent

[Service]
ExecStart=/home/lsjostrom/bin/oidc-agent server

[Install]
WantedBy=default.target
